# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
environment=VNC_PASSWORD="%(ENV_VNC_PASSWORD)s"

[program:xvfb]
command=/usr/bin/Xvfb :99 -screen 0 800x600x16 -ac +extension GLX +extension RANDR
autorestart=true
priority=10
startsecs=2

[program:x11vnc]
command=/docker/start-vnc.sh
environment=VNC_PASSWORD="%(ENV_VNC_PASSWORD)s"
stdout_logfile=/var/log/supervisor/x11vnc.log
stderr_logfile=/var/log/supervisor/x11vnc.log
autorestart=true
priority=20
startsecs=2

[program:flash]
command=/docker/start-flash.sh
environment=DISPLAY=":99"
stdout_logfile=/var/log/supervisor/flash.log
stderr_logfile=/var/log/supervisor/flash.log
autorestart=true
priority=30
startsecs=5

[program:novnc]
command=/bin/bash -c 'sleep 3 && exec /opt/noVNC/utils/novnc_proxy --vnc localhost:5900 --listen 6080 --cert /opt/noVNC/self.pem --ssl-only --heartbeat 30'
autorestart=true
priority=40
startsecs=5
